# deploy/ai-tasks.service
[Unit]
Description=ai-tasks (Next.js task manager)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aitasks
Group=aitasks
WorkingDirectory=/srv/ai-tasks

# Load secrets from the production env file. The path is a symlink that
# resolves to /srv/ai-tasks-data/.env.production — config lives outside
# the deploy dir so it survives wipes.
EnvironmentFile=/srv/ai-tasks/.env.production

Environment=HOSTNAME=127.0.0.1
Environment=NODE_ENV=production

ExecStart=/usr/bin/npm run start --silent

Restart=on-failure
RestartSec=5
TimeoutStopSec=20

# --- Hardening ----------------------------------------------------
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
# Both the deploy dir AND the data dir need write access (data/uploads
# inside the deploy dir are symlinks pointing into /srv/ai-tasks-data).
ReadWritePaths=/srv/ai-tasks /srv/ai-tasks-data

LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-tasks

[Install]
WantedBy=multi-user.target
